# guni_config/face_admin.service
[Unit]
Description=Gunicorn instance to serve core admin service
After=network.target

[Service]
Type=simple
Restart=always
User=pumpk
Group=www-data
WorkingDirectory=/home/pumpk/face-service
Environment="PATH=/home/pumpk/env/bin"
ExecStart=/home/pumpk/env/bin/gunicorn --workers 1 --bind 0.0.0.0:7001 a_wsgi:app

[Install]
WantedBy=multi-user.target